# zfs-killsnaps/debian/zfs-killsnaps.timer
[Unit]
Description=Run ZFS Snapshot Cleanup Weekly

[Timer]
OnCalendar=Sun 02:00
Persistent=true

[Install]
WantedBy=timers.target

